# tools/build.cron
10 0 * * * sh /var/www/html/dawgsquad/daily_build.sh