# tools/test.cron
0 0 * * * /var/www/html/dawgsquad/tools/TEST.pl -a
