# tools/test.cron
# automated test
0 0 * * * perl /var/www/html/dawgsquad/tools/TEST.pl -e
10 0 * * * sh /var/www/html/dawgsquad/tools/daily_build.sh
